# backend_host/systemd/archive-manifest.service
[Unit]
Description=Archive Manifest Generator (24h HLS manifests)
After=network.target

[Service]
Type=simple
User=sunri-pi1
WorkingDirectory=/home/sunri-pi1/virtualpytest/backend_host/scripts
Environment="PATH=/home/sunri-pi1/virtualpytest/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/sunri-pi1/virtualpytest/venv/bin/python3 /home/sunri-pi1/virtualpytest/backend_host/scripts/archive_manifest_generator.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

