# backend_host/systemd/transcript-stream.service
[Unit]
Description=Transcript Accumulator (24h circular audio transcription)
After=network.target

[Service]
Type=simple
User=sunri-pi1
WorkingDirectory=/home/sunri-pi1/virtualpytest/backend_host/scripts
Environment="PATH=/home/sunri-pi1/virtualpytest/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/sunri-pi1/virtualpytest/venv/bin/python3 /home/sunri-pi1/virtualpytest/backend_host/scripts/transcript_accumulator.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

