# deployment/systemd/backend-server.service
[Unit]
Description=VirtualPyTest Backend-Server Service (API Server)
After=network.target
Wants=network.target

[Service]
Type=simple
User=virtualpytest
Group=virtualpytest
WorkingDirectory=/opt/virtualpytest/backend-server
Environment=SERVER_PORT=5109
Environment=SERVER_URL=http://localhost:5109
Environment=DEBUG=false
Environment=CORS_ORIGINS=https://your-frontend.vercel.app

# Git pull and start service
ExecStartPre=/usr/bin/git -C /opt/virtualpytest pull
ExecStart=/usr/bin/python3 src/app.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 