# backend_host/config/services/mac/monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.virtualpytest.monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd ~/virtualpytest/backend_host/scripts && ~/virtualpytest/venv/bin/python capture_monitor.py</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOST_NAME</key>
        <string>mac-host</string>
        <key>PYTHONPATH</key>
        <string>~/virtualpytest/backend_host/scripts</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>StandardOutPath</key>
    <string>/tmp/com.virtualpytest.monitor.out</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/com.virtualpytest.monitor.err</string>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
